# pyLFTPdaemon.conf
# LFTP daemon service

description "LFTP daemon service"
author "SomeOne"

start on runlevel [2345]
stop on runlevel [016]
exec python /opt/pyLFTPdaemon/pyLFTPdaemon.py start
